# ZDNetBackup.service
[Unit]
Description=ZDNetworkBackup
After=network.target

[Service]
User=your_username_here
WorkingDirectory=/full/path/to/your/project
ExecStart=/bin/bash start.sh
Restart=always
RestartSec=5
Environment=PATH=/usr/bin:/usr/local/bin:/bin

[Install]
WantedBy=multi-user.target
